# inventory.service
[Unit]
Description=Inventory Service
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
ExecStartPre=/usr/bin/docker pull quay.io/taylor_fahlman/inventory
ExecStart=/usr/bin/docker run  --rm -d -e HOSTNAME=`hostname -f` --privileged quay.io/taylor_fahlman/inventory

[Install]
WantedBy=multi-user.target
